# syncoid.service
[Unit]
Description=Syncoid Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncoid_backup.sh
