# syncoid.timer
[Unit]
Description=Timer for Syncoid Backups

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
