# syncoid-runner.service
[Unit]
Description=Syncoid Runner Cron Generation Service
Documentation=https://github.com/tgunr/sanoid
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncoid-runner --generate-cron
User=root
IOSchedulingClass=idle
Nice=19

[Install]
WantedBy=multi-user.target
